# systemd/hermes-gateway.service
# /home/YOU/.config/systemd/user/hermes-gateway.service
# Install:  cp hermes-gateway.service ~/.config/systemd/user/
# Enable:   systemctl --user enable hermes-gateway
# Start:    systemctl --user start hermes-gateway
# Logs:     journalctl --user -u hermes-gateway -f

[Unit]
Description=Hermes Agent Gateway
After=network.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
# Adjust the path to your Hermes installation
ExecStart=/home/YOU/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
WorkingDirectory=/home/YOU/.hermes/hermes-agent
Environment="PATH=/home/YOU/.hermes/hermes-agent/venv/bin:/usr/bin:/home/YOU/.local/bin:/usr/local/bin"
Environment="VIRTUAL_ENV=/home/YOU/.hermes/hermes-agent/venv"
Environment="HERMES_HOME=/home/YOU/.hermes"
# — Secrets are loaded from ~/.hermes/.env (chmod 600) — NOT inlined here.
# Inline Environment= values are visible via `systemctl show` and `journalctl`,
# leaking secrets to any local user who can run those commands.
#
# Your ~/.hermes/.env should contain:
#   API_SERVER_KEY=<generated with: python3 -c "import secrets; print(secrets.token_urlsafe(32))">
#   TELEGRAM_BOT_TOKEN=<from @BotFather>
#   TELEGRAM_OWNER_ID=<your numeric Telegram user ID>
#   TELEGRAM_ALLOWED_USERS=<comma-separated allowed user IDs>
#
# Lock down the file after editing:
#   chmod 600 ~/.hermes/.env
EnvironmentFile=%h/.hermes/.env
Restart=on-failure
RestartSec=30
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
